# doc/config/report.ini
[program:report]
command=/opt/report-server/report-startup
directory=/opt/report-server/
priority=60
startsecs=2
autostart=true
autorestart=true
stdout_logfile=/data/logs/report/report-server.stdout.log
stderr_logfile=/data/logs/report/report-server.stderr.log
